# roles/ps1auth/templates/ps1auth-web.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream={{ app_root }}/run/gunicorn.sock

[Install]
WantedBy=sockets.target
